# misc/xdg-desktop-portal-shana.service.in
[Unit]
Description=Xdg Desktop Portal for other portal
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.shana
ExecStart=@xdg_install_dir@/xdg-desktop-portal-shana
Restart=on-failure
